# app/scripts/chess-backend.service
[Unit]
Description=Turniej Szachowy — Backend FastAPI
After=network.target postgresql.service

[Service]
User=chess
WorkingDirectory=/opt/chess/backend
EnvironmentFile=/opt/chess/backend/.env
ExecStart=/opt/chess/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
